# deploy/report_server.conf
[program:market_monitor]
command=/usr/bin/php /home/ubuntu/CryptoArbitrage/ReportingServer.php --mongodb --monitor=0
autostart=true
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog